# docker-openshift/docker/s6/openssh/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

echo "test"

exec gosu root /usr/sbin/sshd -D -f /app/gogs/docker/sshd_config
